# files/services.conf
[program:nginx]
command = /usr/sbin/nginx
user = root
autostart = true
[program:php5-fpm]
command = /usr/sbin/php-fpm
user = root
autostart = true
